# roles/rclone/files/rclone-sync.service
[Unit]
Description=rclone-sync

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone-sync
